# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:backend]
command=gunicorn arayuz:app --bind 127.0.0.1:5000 --workers 2 --threads 4 --timeout 60
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=sh -c "cd zedin-frontend && npx next start --port %(ENV_PORT)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

